# firewall-agent.service
[Unit]
Description=Whitelist Firewall Agent (ipset + iptables)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/root/shild-captive-portal
ExecStart=/usr/bin/python3 /root/shild-captive-portal/firewall_agent.py
EnvironmentFile=/root/shild-captive-portal/.env
Restart=always
RestartSec=10
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
